# scripts/aima.service
[Unit]
Description=AI Inference Managed by AI
After=network.target

[Service]
Type=simple
User=aima
Group=aima
ExecStart=/usr/local/bin/aima start
Restart=on-failure
RestartSec=5

LimitNOFILE=65536

Environment="AIMA_CONFIG=/etc/aima/aima.toml"

[Install]
WantedBy=multi-user.target
